# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:warphish]
command=python3 /app/warphish.py
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/warphish/error.log
stdout_logfile=/var/log/warphish/out.log
user=root
environment=PYTHONUNBUFFERED="1"

[program:web-dashboard]
command=gunicorn --bind 0.0.0.0:5000 --workers 4 --threads 2 --timeout 60 warphish:app
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/warphish/web-error.log
stdout_logfile=/var/log/warphish/web-out.log

[program:redis]
command=redis-server --appendonly yes
autostart=true
autorestart=true
stderr_logfile=/var/log/redis/redis-error.log
stdout_logfile=/var/log/redis/redis-out.log

[program:postgres]
command=/usr/local/bin/docker-entrypoint.sh postgres
autostart=true
autorestart=true
stderr_logfile=/var/log/postgres/error.log
stdout_logfile=/var/log/postgres/out.log